# bsp/overrides/bluetooth.service.d/override.conf
[Service]
# Blank ExecStart line to clear the service file, overrides are additive.
ExecStart=
ExecStart=/usr/libexec/bluetooth/bluetoothd --noplugin=sap --experimental
